# misc/root.cron
43 6 10 * * certbot renew --post-hook "systemctl reload nginx"